# tar1090-analytics-ingest.service
[Unit]
Description=tar1090 analytics ingest (aircraft.json to database)
After=network.target postgresql.service tar1090-analytics-api.service
Wants=postgresql.service

[Service]
User=tar1090
Group=tar1090
WorkingDirectory=/usr/local/share/tar1090/analytics/ingest
EnvironmentFile=/etc/default/tar1090-analytics
ExecStart=/usr/local/share/tar1090/analytics-venv/bin/python main.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
